# crontab/crontab.sh
5 8 * * * cd /Users/admin && ./git_update.sh
10 8 * * * cd /Users/admin && ./tar_printers.sh
15 8 * * * cd /Users/admin && ./git_shorten.sh
