# crontab/git_shorten.sh
#!/bin/sh
PATH=/usr/local/bin:/usr/local/sbin:~/bin:/usr/bin:/bin:/usr/sbin:/sbin

cp -p /Users/admin/github/*.sh /Users/admin/smallhub
cp -p /Users/admin/github/*.bat /Users/admin/smallhub
cp -Rp /Users/admin/github/WITup\ app /Users/admin/smallhub
